# srcs/grafana/supervisor.conf
[supervisord]
nodaemon=true

[program:grafana]
command=sh -c "cd /grafana-7.5.1/bin/ && ./grafana-server && kill -s SIGTERM $(cat supervisord.pid)"

# [program:telegraf]
# command=sh -c "telegraf -config /etc/telegraf.conf && kill -s SIGTERM $(cat supervisord.pid)"